# system_files/etc/systemd/system/machine-id.service
# Ensure local appliance machine identity exists before agent starts.
[Unit]
Description=Initialize appliance machine identity
Documentation=man:systemd.service(5)
Before=agent.service

[Service]
Type=oneshot
ExecStart=/usr/lib/appcoreos/init-machine-id.sh
RemainAfterExit=yes

[Install]
WantedBy=appcore.target
