# system_files/etc/systemd/system/tui.service
[Unit]
Description=AppCoreOS TUI
After=systemd-user-sessions.service
Conflicts=getty@tty1.service serial-getty@ttyS0.service console-getty.service

[Service]
ExecStart=/usr/lib/appcoreos/tui.sh
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/ttyS0
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
