# rootfs/etc/systemd/system/appcore.target
[Unit]
Description=AppCoreOS Target
Requires=tui.service
After=basic.target

[Install]
Alias=default.target
