# rootfs/etc/systemd/system/getty@ttyS0.service
[Unit]
Description=Disabled getty on ttyS0

[Service]
ExecStart=
ExecStart=/bin/false
