# rootfs/etc/systemd/system/getty@ttyS0.service
[Unit]
Description=Disabled console login

[Service]
ExecStart=
ExecStart=/bin/false
