# rootfs/etc/systemd/system/tui.service
[Unit]
Description=AppCoreOS TUI
After=systemd-user-sessions.service

[Service]
ExecStartPre=-/usr/sbin/setenforce 0
ExecStart=/usr/lib/your-os/tui.sh
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/ttyS0
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
Restart=always

[Install]
WantedBy=multi-user.target
