# rootfs/etc/systemd/system/tui.service
[Unit]
Description=AppCoreOS TUI
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
ExecStart=/usr/lib/your-os/tui.sh
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
Restart=always

[Install]
WantedBy=default.target
