# rootfs/etc/systemd/system/update-os.service
# Run scheduled OS self-update workflow.
[Unit]
Description=Run OS self-update and reboot
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/your-os/update-os.sh
